# bee-tracker.service
[Unit]
Description=Bee Tracker — YOLO + BoT-SORT
After=network.target

[Service]
ExecStart=/home/pi/beeTracker/venv/bin/python /home/pi/beeTracker/live_tracker.py --no-display
WorkingDirectory=/home/pi/beeTracker
Restart=always
RestartSec=10
User=pi
Environment=DISPLAY=

[Install]
WantedBy=multi-user.target
